# .config/systemd/user/mpd-notify.service
[Unit]
Description=Notification daemon for MPD

[Service]
ExecStart=/home/dophin2009/.local/bin/mpd-notify.go

[Install]
WantedBy=default.target
