# .config/systemd/user/mpdnotify.service
[Unit]
Description=Notification daemon for MPD

[Service]
ExecStart=/home/dophin2009/.local/bin/mpdnotify -d

[Install]
WantedBy=default.target
