# battery-notify/.config/systemd/battery-monitor.service
[Unit]
Description=Monitor battery status

[Service]
ExecStart=/home/dophin2009/.local/bin/monitor-battery

[Install]
WantedBy=default.target
