# .config/sv/mpdnotify/run
#!/bin/sh
exec 2>&1
exec $HOME/.local/bin/mpdnotify monitor
