# .config/sv/syncthing/run
#!/bin/sh
export HOME=/home/dophin2009
exec 2>&1
exec syncthing 
