# packages/battery-notify/tree/.config/systemd/user/battery-monitor.service.tmpl
[Unit]
Description=Monitor battery status

[Service]
ExecStart={{ .home }}/.local/bin/monitor-battery

[Install]
WantedBy=default.target
