# hosts/segv/sys/init.d/hydroxide
#!/sbin/openrc-run
command="hydroxide"
command_args="serve"
command_background=true
pidfile="/run/hydroxide.pid"

command_user="${USER}"
