# hosts/segv/sys/init.d/hydroxide
#!/sbin/openrc-run
command="hydroxide"
command_args="serve"

command_user="${USER}"
