# hosts/segv/sys/init.d/protonmail-bridge
#!/sbin/openrc-run
command="protonmail-bridge"
command_args="--noninteractive"
command_background=true
pidfile="/run/protonmail-bridge.pid"

command_user="${USER}"
