# mpdnotify/systemd-service/tree/.config/systemd/user/mpdnotify.service.tmpl
[Unit]
Description=Notification daemon for MPD
After=mpd.service

[Service]
ExecStart={{ .bin }} monitor

[Install]
WantedBy=default.target
